# 资料文档/supervisordphp.ini
[program:php]
command=php -f /www/path/ab.php
autorstart=true
autorestart=true
stdout_logfile=/tmp/supervisor.log
